# games-site/auto-build-site.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.games-site.auto-build-site</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/zhaoq0103/DEV/Personal/How-to-Make-Money/games-site/auto-build-site.sh</string>
    </array>

    <!--
      每天运行两次，以美国东部时间 (ET = UTC-4 夏令时) 为准：

      🌅 Morning run: 10 AM ET  = 14:00 UTC = 22:00 China
         覆盖美国用户早间搜索高峰 + Steam 日间内容更新

      🌙 Evening run: 7 PM ET   = 23:00 UTC = 07:00 China (次日)
         覆盖美国晚间游戏高峰（下班/放学后打游戏、搜攻略）
    -->

    <key>StartCalendarInterval</key>
    <array>
        <!-- Morning: 每天 22:07 (北京时间) = 10:07 AM ET -->
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>7</integer>
        </dict>
        <!-- Evening: 每天 07:07 (北京时间) = 7:07 PM ET -->
        <dict>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>7</integer>
        </dict>
    </array>

    <!-- 日志 -->
    <key>StandardOutPath</key>
    <string>/Users/zhaoq0103/Library/Logs/auto-build-site/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhaoq0103/Library/Logs/auto-build-site/stderr.log</string>

    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/zhaoq0103</string>
        <key>PATH</key>
        <string>/Users/zhaoq0103/.local/bin:/Users/zhaoq0103/.nvm/versions/node/v24.4.0/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
